# ansible/roles/ansible_pull/templates/ansible-pull.service.j2
[Unit]
Description=Ansible Pull

[Service]
Type=oneshot
ExecStart=/bin/bash {{ pull_home }}/ansible/scripts/pull.sh update.yml