# ansible/roles/ansible_pull/templates/ansible-pull.timer.j2
[Unit]
Description=Run Ansible Pull Periodically

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target