# ansible/roles/ansible_pull/templates/ansible-pull.service.j2
[Unit]
Description=Ansible Pull

[Service]
Type=oneshot
ExecStart={{ pull_home }}/scripts/ansible/pull.sh update.yml